# playwright-mcp.service
[Unit]
Description=Playwright MCP Server
Documentation=https://github.com/modelcontextprotocol/servers
After=network.target
Wants=network.target

[Service]
Type=simple
User=mauri
Group=mauri
WorkingDirectory=/home/mauri/mcp-agent
ExecStart=/usr/local/bin/npx @playwright/mcp@latest --host 127.0.0.1 --port 8931 --user-data-dir /home/mauri/mcp-agent/.pw-profile --shared-browser-context --caps vision --console-level debug --snapshot-mode full --timeout-action 45000 --timeout-navigation 90000 --viewport-size 1280x720 --init-page /home/mauri/mcp-agent/src/mcp/init-linkedin-page.ts --allowed-hosts "*"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=playwright-mcp

# Environment
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/mauri

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/mauri/mcp-agent
ProtectHome=read-only

# Resource limits
LimitNOFILE=65536
MemoryHigh=2G
MemoryMax=4G

[Install]
WantedBy=multi-user.target